# services/power-manager/power-manager.service
[Unit]
Description=CarStereo Ignition Power Manager
Documentation=https://github.com/kdoehler/carStereo
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/carstereo/services/power-manager/power-manager.py --grace-seconds 30 --usb-park
Restart=on-failure
RestartSec=10

# Needs root for GPIO and CPU governor writes
User=root

# Logging
StandardOutput=append:/var/log/carstereo/power-manager.log
StandardError=append:/var/log/carstereo/power-manager.log

# Hardening
# NOTE: ProtectSystem=strict is NOT used here. It blocks waydroid's
# LXC status checks (/run, /proc, /sys paths), causing "waydroid
# container freeze" to silently fail. An unfrozen waydroid during
# USB power cycling kills touchscreen re-enumeration on wake.
ProtectHome=read-only
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
